# addon/imou-p2p-bridge/rootfs/etc/services.d/imou-p2p-bridge/run
#!/bin/sh
exec python3 /opt/imou-p2p-bridge/supervisor.py
